# 11-usbdisconnect.rules
SUBSYSTEM=="usb", ACTION=="unbind", RUN+="/usr/bin/usbdevgone.sh $devpath"
